# docker_symfony/php_fpm/crontab
@daily /usr/local/bin/php /var/www/bin/console parse:run
